# crontab
*/5 * * * * * *  timeout 2 python ./job.py
